# scripts/launchd/com.claudemind.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.claudemind.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>--import</string>
    <string>tsx</string>
    <string>.claude/bin/lkhs-daemon.ts</string>
  </array>
  <key>WorkingDirectory</key><string>__VAULT__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>__LOGDIR__/launchd.daemon.out.log</string>
  <key>StandardErrorPath</key><string>__LOGDIR__/launchd.daemon.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__PATH__</string>
    <key>LKHS_VAULT_ROOT</key><string>__VAULT__</string>
  </dict>
</dict>
</plist>
