# scripts/launchd/com.claudemind.dream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.claudemind.dream</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>--import</string>
    <string>tsx</string>
    <string>.claude/bin/pipelines.ts</string>
    <string>auto_dream</string>
  </array>
  <key>WorkingDirectory</key><string>__VAULT__</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>__LOGDIR__/launchd.dream.out.log</string>
  <key>StandardErrorPath</key><string>__LOGDIR__/launchd.dream.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__PATH__</string>
    <key>LKHS_VAULT_ROOT</key><string>__VAULT__</string>
  </dict>
</dict>
</plist>
